# claudio.service
[Unit]
Description=Cláudio v2 — assistente pessoal
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/share/claudio-venv/bin/python -m claudio.main
WorkingDirectory=%h/repos/projetos/claudio-v2
Environment=PYTHONPATH=%h/repos/projetos/claudio-v2
Environment=PYTHONUNBUFFERED=1
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
